# prod/systemd/certbot-renew.timer
[Unit]
Description=Runs certbot-renew.service daily

[Timer]
OnBootSec=30d
OnUnitActiveSec=30d
Persistent=true

[Install]
WantedBy=timers.target
